# installer/systemd/ee-agent.service
[Unit]
Description=Easy Enclave Deployment Agent
After=network.target libvirtd.service
Wants=libvirtd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/easy-enclave
EnvironmentFile=/etc/easy-enclave/agent.env
ExecStart=/opt/easy-enclave/venv/bin/python /opt/easy-enclave/agent.py
Restart=always
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1

# Security hardening (allow libvirt access)
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=false
# PrivateTmp must be false so libvirt can access the temp files
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ee-agent

[Install]
WantedBy=multi-user.target
