# agent/systemd/ee-agent.service
[Unit]
Description=Easy Enclave Deployment Agent
After=network.target libvirtd.service
Wants=libvirtd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/easy-enclave
ExecStart=/usr/bin/python3 /opt/easy-enclave/agent.py --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1

# Security hardening (allow libvirt access)
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=false
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ee-agent

[Install]
WantedBy=multi-user.target
