# installer/templates/agent-service.service
[Unit]
Description=Easy Enclave Agent (VM)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
Environment=EE_AGENT_MODE=vm
Environment=EE_AGENT_DIR=/opt/ee-agent
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/easy-enclave/agent.env
WorkingDirectory=/opt/ee-agent
ExecStart=/opt/ee-agent/venv/bin/python /opt/ee-agent/agent.py
Restart=always
RestartSec=2
StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=ee-agent

[Install]
WantedBy=multi-user.target
